# etc/systemd/system/docker.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -D \
	--storage-driver overlay2
	-H fd://
